# guacamole/supervisor-proxy.ini
[program:proxy]
command=/opt/proxy-venv/bin/python /opt/proxy.py
autostart=true
autorestart=true
priority=950
environment=PYTHONUNBUFFERED="1"
stdout_logfile=/var/log/supervisor/proxy.log
stderr_logfile=/var/log/supervisor/proxy_err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
